# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
nohup node /opt/jd/panel/server.js >/dev/null 2>&1 &
exit 0
